# examples/vacon-signaling-server.socket
[Unit]
Description=Vacon signaling server socket

[Socket]
ListenStream=/run/vacon-signaling-server.websocket

[Install]
WantedBy=sockets.target
